# templates/carvisor.service.j2
[Unit]
Description=Carvisor service
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/python /usr/src/carvisor/main.py

[Install]
WantedBy=multi-user.target
